# conf/supervisor/rq-dashboard.conf
[program:rq-dashboard]
command=%(ENV_HOME)s/.pyenv/shims/rq-dashboard

stdout_logfile=%(ENV_HOME)s/tmp/rq-dashboard.log ; stdout log path
stdout_logfile_maxbytes=5MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)

autostart=true
autorestart=true
stopsignal=TERM
