# conf/supervisor/stack_oprc_job.conf
[program:stack_oprc_job]
;
; run ~/lab/src_lib/python/stack_oprc_job.py like deamon
; you need
; $ chmod +x ~/lab/src_lib/python/stack_oprc_job.py
;
command = %(ENV_HOME)s/lab/src_lib/python/stack_oprc_job.py

; stdout log path
stdout_logfile=%(ENV_HOME)s/tmp/%(ENV_HOST)s_stack_oprc_job.log
stdout_logfile_maxbytes=20MB ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10    ; # of stdout logfile backups (default 10)

autostart=true
autorestart=true
stopsignal=TERM
