# config/autostart/set_t480s.desktop
[Desktop Entry]
Name=t460s_settings
Type=Application
Exec=/home/yyamnk/bin/set_t480s.sh
X-MATE-Autostart-enabled=true
