# config/autostart/set_t480s.desktop
[Desktop Entry]
Type=Application
Exec=/home/bin/set_t480s.sh
Hidden=false
X-MATE-Autostart-enabled=true
Name=t460s_settings
